# touch.service
[Unit]
Description=Touch Driver

[Service]
ExecStart=/etc/init.d/touch.sh

[Install]
WantedBy=multi-user.target